# templates/helga.ini.j2
[program:helga]
directory={{ helga_path }}
command=helga --settings=settings.py
user=helga
autostart=true
autorestart=true
redirect_stderr=True
